# etc/supervisor/conf.d/logstash.conf
[program:logstash]
command=/usr/share/logstash/bin/logstash -f /etc/logstash/conf.d/logstash.conf
stdout_logfile=/var/log/logstash/logstash.log
redirect_stderr=true
